# sources/supervisor/rsyncd.ini
[program:rsyncd]
command=/usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf
priority=999
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
