# sources/supervisor/tftpd.ini
[program:tftpd]
command=/usr/sbin/in.tftpd -s /var/lib/tftpboot -L -v
priority=999
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
