# sources/supervisor/uwsgi.ini
[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi.d/cobbler-svc.ini --ini /etc/uwsgi.d/cobbler-web.ini
priority=999
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
